# dundee/dundee.service.in
[Unit]
Description=DUN service
After=syslog.target

[Service]
Type=dbus
BusName=org.ofono.dundee
ExecStart=@prefix@/sbin/dundee -n
StandardError=null

[Install]
WantedBy=multi-user.target
